# packages/linux/tor/etc/systemd/system/tor.service.d/override.conf
[Service]
BindReadOnlyPaths=/home/justin/workspace/dotfiles/tor:/home/justin/workspace/dotfiles/tor
ProtectHome=tmpfs
